# etc/init.d/bismark
#!/bin/sh /etc/rc.common

START=90
STOP=1

boot() {
	crontab -l | grep -vi "# remove on reboot" | crontab -
	bismark-bootstrap
}
